# scripts/com.kelok.afcalendar.sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Weekly AF Calendar timetable sync.
  Fires every Sunday at 22:30 (your Mac's local time — set the Mac to MYT so this
  lands after the Sunday 22:00 MYT cache rollover). If the Mac is asleep/off at
  22:30, launchd runs it at the next wake. Logs to sync.log in the project dir.

  Install:
    cp scripts/com.kelok.afcalendar.sync.plist ~/Library/LaunchAgents/
    launchctl load ~/Library/LaunchAgents/com.kelok.afcalendar.sync.plist
  Run immediately to test:
    launchctl start com.kelok.afcalendar.sync
  Remove:
    launchctl unload ~/Library/LaunchAgents/com.kelok.afcalendar.sync.plist
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.kelok.afcalendar.sync</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/kelokchan/Work/af-calendar/scripts/run-sync.sh</string>
  </array>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Weekday</key>
    <integer>0</integer>
    <key>Hour</key>
    <integer>22</integer>
    <key>Minute</key>
    <integer>30</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/kelokchan/Work/af-calendar/sync.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/kelokchan/Work/af-calendar/sync.log</string>

  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
